# incident-logger.service
[Unit]
Description=AI Incident Logger API
After=network.target

[Service]
User=claw
WorkingDirectory=/home/claw/Projects/ai-incident-logger
EnvironmentFile=/home/claw/Projects/ai-incident-logger/.env
ExecStart=/home/claw/Projects/ai-incident-logger/venv/bin/gunicorn api:app --bind 127.0.0.1:5001 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
